# meters.service
[Unit]
Description=meters

[Service]
ExecStart=/home/pi/meters.py

[Install]
WantedBy=multi-user.target
